# home/dot_config/autostart/1password.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/1password --silent
Comment=1Password
Name=1Password
