# home/dot_config/autostart/beeper.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/beeper
Comment=Beeper
Name=Beeper
