# home/dot_config/autostart/cava.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/cava
Comment=Cava
Name=Cava
