# home/dot_config/autostart/discord.desktop
[Desktop Entry]
Type=Application
Exec=/bin/flatpak run com.discordapp.Discord
Comment=Discord
Name=Discord
