# home/dot_config/autostart/discord.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/discord
Comment=Discord
Name=Discord
