# home/dot_config/autostart/element.desktop
[Desktop Entry]
Type=Application
Exec=/bin/flatpak run im.riot.Riot
Comment=Element
Name=Element
