# home/dot_config/autostart/element.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/element-desktop
Comment=Element
Name=Element
