# home/dot_config/autostart/firefox-developer-edition.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/firefox-developer-edition
Comment=Firefox
Name=Firefox
