# home/dot_config/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Exec=/bin/flatpak run org.mozilla.firefox
Comment=Firefox
Name=Firefox
