# home/dot_config/autostart/firefox-developer-edition-personal.desktop
[Desktop Entry]
Type=Application
Exec=/usr/lib/firefox-developer-edition/firefox --class="firefox-personal" -P personal
Comment=Firefox
Name=Firefox
