# home/dot_config/autostart/librewolf.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/librewolf
Comment=Librewolf
Name=Librewolf
