# home/dot_config/autostart/mako.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/mako
Comment=Notifications daemon
Name=Mako
