# home/dot_config/autostart/slack.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/slack
Comment=Slack
Name=Slack
