# home/dot_config/autostart/slack.desktop
[Desktop Entry]
Type=Application
Exec=/bin/flatpak run com.slack.Slack
Comment=Slack
Name=Slack
