# home/dot_config/systemd/user/bugwarrior-pull.service.tmpl
{{ if .taskwarrior }}
[Unit]
Description=bugwarrior-pull

[Service]
Type=oneshot
ExecStart=/home/joshmock/.local/bin/bugwarrior-pull

[Install]
WantedBy=bugwarrior-pull.timer
{{ end }}
