# home/dot_config/systemd/user/bugwarrior-pull.timer.tmpl
{{ if .taskwarrior }}
[Unit]
Description=bugwarrior-pull (hourly)

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
{{ end }}
