# home/dot_config/systemd/user/calendar-reminder.timer.tmpl
{{ if .work }}
[Unit]
Description=Calendar reminders

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
{{ end }}
