# home/dot_config/systemd/user/calendar-reminder.service.tmpl
{{ if .work }}
[Unit]
Description=Calendar reminders

[Service]
Type=oneshot
ExecStart=/home/joshmock/.local/bin/calendar-notifier.sh

[Install]
WantedBy=calendar-reminder.timer
{{ end }}
