# home/dot_config/systemd/user/clear-mail-archive.timer.tmpl
{{ if .email -}}
[Unit]
Description=remove oldest messages from email archive

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
{{ end }}
