# home/dot_config/systemd/user/empty-mail-trash.timer.tmpl
{{ if .email -}}
[Unit]
Description=empty email trash

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
{{ end }}
