# home/dot_config/systemd/user/mopidy-queue.timer.tmpl
{{- if .mopidy }}
[Unit]
Description=Keep the music playing

[Timer]
OnCalendar=*-*-* *:00,05,10,15,20,25,30,35,40,45,50,55:00

[Install]
WantedBy=timers.target
{{ end }}
