# home/dot_config/systemd/user/mopidy-queue.service.tmpl
{{- if .mopidy }}
[Unit]
Description=Keep the music playing

[Service]
Type=oneshot
ExecStart=/home/joshmock/.local/bin/mopidy-queue-manager

[Install]
WantedBy=mopidy-queue.timer
{{ end }}
