# home/dot_config/systemd/user/purge-old-tasks.timer.tmpl
{{ if .taskwarrior -}}
[Unit]
Description=purge old tasks

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
{{ end }}
